# crontabs/track_my_ip
#!/bin/bash

echo "$(date) $(getip) $(ipub)" >> ~/my_ip_logs


